# systemd/watchgoat.service
[Unit]
Description=Edgise WatchGoat

[Service]
Type=oneshot
ExecStart=PATHTOEXECUTABLE