# systemd/watchgoat.service
[Unit]
Description=WatchGoat

[Service]
Type=oneshot
ExecStart=PATHTOEXECUTABLE